# supervisord.conf
[supervisord]
nodaemon=false
logfile=/var/log/supervisor/supervisord.log
logfile_maxbytes=10MB
loglevel=info
pidfile=/var/run/supervisord.pid

; ── xrdp session manager ─────────────────────────────────────────────────────
[program:xrdp-sesman]
command=/usr/sbin/xrdp-sesman --nodaemon
autostart=true
autorestart=true
priority=20
stdout_logfile=/var/log/supervisor/xrdp-sesman.log
stderr_logfile=/var/log/supervisor/xrdp-sesman.log

; ── xrdp (RDP server, listens on port 3389) ──────────────────────────────────
[program:xrdp]
command=/usr/sbin/xrdp --nodaemon
autostart=true
autorestart=true
priority=30
stdout_logfile=/var/log/supervisor/xrdp.log
stderr_logfile=/var/log/supervisor/xrdp.log

; ── openclaw gateway ─────────────────────────────────────────────────────────
; Runs as the 'node' user, listens on port 18789 (web UI) and 18790 (bridge).
[program:openclaw-gateway]
command=node /app/dist/index.js gateway --bind lan --port 18789
user=node
environment=HOME="/home/node",TERM="xterm-256color",OPENCLAW_GATEWAY_TOKEN="%(ENV_OPENCLAW_GATEWAY_TOKEN)s"
directory=/app
autostart=true
autorestart=true
priority=40
stopwaitsecs=10
stdout_logfile=/var/log/supervisor/openclaw-gateway.log
stderr_logfile=/var/log/supervisor/openclaw-gateway.log
